# ULWGL/sniper_platform_0.20240125.75305/files/lib/systemd/user/glib-pacrunner.service
[Unit]
Description=GLib proxy auto-configuration service

[Service]
Type=dbus
BusName=org.gtk.GLib.PACRunner
ExecStart=/usr/libexec/glib-pacrunner
